# ansible/roles/mailpilot_crm/templates/mailpilot.service.j2
[Unit]
Description=MailPilot agent service
After=postgresql.service network-online.target
Wants=network-online.target
# Schema-gate / config-class failures exit 1. Restart=on-failure alone
# crash-loops and systemctl can report active while the agent never stays
# up. Start-rate limit forces a durable failed state operators can alert on.
# RestartPreventExitStatus not used yet — app has no dedicated config exit code.
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User={{ mailpilot_user }}
Group={{ mailpilot_group }}
WorkingDirectory={{ mailpilot_home }}
{% if google_application_credentials_b64 | default('') | length > 0 %}
# Workload Identity is unavailable off-GCP; pin the SA JSON for ADC clients.
Environment=GOOGLE_APPLICATION_CREDENTIALS={{ google_application_credentials_path }}
{% endif %}
ExecStart={{ mailpilot_home }}/.local/bin/mailpilot run
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
